# scripts/launchd/com.benlangsfeld.residence-exec-prep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.benlangsfeld.residence-exec-prep</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$HOME/Documents/Apps/grain" &amp;&amp; /opt/homebrew/bin/npx tsx scripts/residence-exec-prep.ts</string>
    </array>

    <!-- Mondays at 6:33am local. Off-zero minute. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>6</integer>
        <key>Minute</key><integer>33</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ben/Library/Logs/residence-exec-prep.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ben/Library/Logs/residence-exec-prep.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
